# admin/kausapp-admin.service
[Unit]
Description=KausApp Bug-Report Admin (Tailscale-only)
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=acronix
WorkingDirectory=/home/acronix/KausApp/admin
EnvironmentFile=/home/acronix/KausApp/admin/kausapp-admin.env
ExecStart=/usr/bin/python3 /home/acronix/KausApp/admin/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
